# gems/rake-support/share/init/torquebox.conf.erb
description "This is an upstart job file for torquebox"

pre-start script
bash << "EOF"
  mkdir -p /var/log/torquebox
  chown -R torquebox /var/log/torquebox
EOF
end script

start on started network-services
stop on stopped network-services
respawn

exec su - torquebox -c '/opt/torquebox/jboss/bin/standalone.sh <%= server_opts %> >> /var/log/torquebox/torquebox.log 2>&1'
